# s6-overlay/s6-rc.d/rpcbind/run
#!/command/execlineb -P
/usr/sbin/rpcbind -f 
